# supervisor/zookeeper.conf
[program:zookeeper]
command=/usr/bin/start-zookeeper.sh start-foreground
autostart=true
autorestart=true